# deploy_tools/gunicorn.systemd.template.service
[Unit]
Description=Gunicorn server for DOMAIN

[Service]
Restart=on-failure
User=arno
WorkingDirectory=/home/arno/sites/DOMAIN
EnvironmentFile=/home/arno/sites/DOMAIN/.env

ExecStart=/home/arno/sites/DOMAIN/virtualenv/bin/gunicorn \
    --bind unix:/tmp/DOMAIN.socket \
    superlists.wsgi:application

[Install]
WantedBy=multi-user.target
